# wifi_hotspot/rootfs/etc/s6-overlay/s6-rc.d/hostapd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Run hostapd (the access point daemon)
# ==============================================================================
bashio::log.info "Starting hostapd..."

exec hostapd /run/hotspot/hostapd.conf
